# templates/backup-db.cron.j2
SHELL=/bin/sh
MAILTO=root

55 */6 * * * root (nice /usr/local/sbin/backup-db)
